# crontab
49 * * * * jes cd /home/jes/weather && ./run
