# crontab
49 * * * * jes cd /home/jes/motorcycling-forecast && ./run
